# cron.hourly/statistics
#!/bin/sh

# Gateway Statistiken
/usr/local/bin/vnstati.sh
